# nav-checker-api.service
# Systemd service for NAV Checker API
# Place this in /etc/systemd/system/nav-checker-api.service

[Unit]
Description=NAV Checker API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/nav_assistant_project/python_server
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/root/nav_assistant_project/python_server
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8001 --reload
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
